# general/.config/crons.cron
DISPLAY=:0
# 0 23 * * * source $HOME/.Xdbus; dunstify -u 2 "sleep time"
@reboot source $HOME/.config/.Xdbus; cd $HOME/Projects/app; pipenv run python app.py
@monthly cd $HOME/Projects/updateReminder; python main.py
*/1 * * * * source $HOME/.config/.Xdbus; $HOME/bin/battery_low_alert.sh
*/30 * * * * rclone sync ./org remote:org --include *.org
*/30 * * * * cd $HOME/org; $HOME/bin/git-sync


# vim: syntax=crontab
