# root/etc/NetworkManager/dispatcher.d/90-using-box-refresh
#!/bin/sh
/usr/bin/sing-box-refresh-iface
