# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=udiskie automounter for removable media daemon

[Service]
Type=simple
ExecStart=/bin/udiskie

[Install]
WantedBy=default.target
