# dunst/.config/systemd/user/dunst.service
[Unit]
Description=dunst daemon

[Service]
Type=simple
ExecStart=/bin/dunst

[Install]
WantedBy=default.target
