# isync/.config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=forking
ExecStart=/bin/bash -c '. /home/yasser/.config/neomutt/bin/syncmail'

[Install]
WantedBy=graphical.target
