# isync/.config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service
OnFailure=notify-popup-send@%i.service

[Service]
Type=forking
ExecStart=bash -ic "source /home/yasser/.profile && /home/yasser/.config/neomutt/bin/syncmail"

[Install]
WantedBy=graphical.target
