# isync/.config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=5m
Unit=mbsync.service

[Install]
WantedBy=timers.target
