# .config/crons.cron
DISPLAY=:0
0 23 * * * source $HOME/.Xdbus; dunstify "sleep time"
@monthly cd $HOME/Projects/updateReminder; python main.py
*/1 * * * * source $HOME/.config/.Xdbus; $HOME/bin/battery_low_alert.sh

*/30 * * * * rclone sync -u ./org remote:org --include *.org
*/30 * * * * rclone sync -u remote:org ./org --include *.org
*/30 * * * * cd $HOME/org; $HOME/bin/git-sync


# vim: syntax=crontab
