# cron/.config/cron/crons.cron
*/1 * * * * . $HOME/.config/cron/cron_env && $HOME/bin/battery_low_alert

# Sync mega material
0 12 * * * $HOME/bin/wait_internet && rclone sync $HOME/Documents mega:documents

# Sync notes
*/30 * * * * . $HOME/.config/cron/cron_env && bash -ic syncorg && cd $HOME/notes && $HOME/bin/git-sync -n -s

* */1 * * * . $HOME/.config/cron/cron_env && $HOME/.dotfiles-private && $HOME/bin/git-sync -n -s

# vim: syntax=crontab commentstring=#%s
